# deploy/backup/medbot-backup.service
# Системный юнит (раскатывается CI/setup.sh в /etc/systemd/system/).
# Активируется таймером medbot-backup.timer (ежедневно).

[Unit]
Description=Med Bot — PostgreSQL backup

[Service]
Type=oneshot
EnvironmentFile=/root/after-30-medicine-bot/.env
ExecStart=/bin/bash /root/after-30-medicine-bot/deploy/backup/backup.sh
StandardOutput=journal
StandardError=journal
