# deploy/systemd/medbot-api.service
[Unit]
Description=MedBot FastAPI
After=network-online.target postgresql@16-main.service
Wants=network-online.target
Requires=postgresql@16-main.service
# OP1: лимит на crash-loop — 5 рестартов за 5 мин, иначе failed
StartLimitIntervalSec=300
StartLimitBurst=5
# OP3: алерт админу при переходе в failed
OnFailure=medbot-alert@%n.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/after-30-medicine-bot
EnvironmentFile=/root/after-30-medicine-bot/.env
ExecStart=/root/after-30-medicine-bot/venv/bin/uvicorn api.main:app --port 8000
# OP6: sandbox-харднинг. App живёт в /root → выделенного пользователя не вводим
# (потребовал бы релокации каталога), а ограничиваем доступ через ReadWritePaths.
# /admin/stats дёргает systemctl is-active — @system-service разрешает execve.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/root/after-30-medicine-bot /run/postgresql
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
# OP1: always (а не on-failure) — переживает TERM/чистый выход
Restart=always
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
